# luci-app-phantun/root/etc/init.d/phantun
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

get_config() {
	config_get_bool enabled "$1" enabled 1
	config_get_bool ipv4_only "$1" ipv4_only 0
	config_get type "$1" type "client"
	config_get local_addr "$1" local_addr "127.0.0.1"
	config_get local_port "$1" local_port
	config_get remote_addr "$1" remote_addr
	config_get remote_port "$1" remote_port
	config_get tun_name "$1" tun_name
	config_get tun_local "$1" tun_local
	config_get tun_peer "$1" tun_peer
	config_get tun_local6 "$1" tun_local6
	config_get tun_peer6 "$1" tun_peer6
	config_get handshake_file "$1" handshake_file
}

start_instance() {
	local cfg="$1"
	get_config "$cfg"

	[ "$enabled" -eq 1 ] || return

	# SAFETY CHECK: Prevent blocking SSH Access on Router
	if [ "$local_port" = "22" ]; then
		logger -t phantun "CRITICAL: Cannot use port 22 for Phantun local port! Skipping instance $cfg."
		return
	fi

	procd_open_instance

	if [ "$type" = "server" ]; then
		procd_set_param command /usr/bin/phantun_server
		procd_append_param command --local "$local_port"
	else
		procd_set_param command /usr/bin/phantun_client
		procd_append_param command --local "$local_addr:$local_port"
	fi

	procd_append_param command --remote "$remote_addr:$remote_port"

	[ -n "$tun_name" ] && procd_append_param command --tun "$tun_name"
	[ -n "$tun_local" ] && procd_append_param command --tun-local "$tun_local"
	[ -n "$tun_peer" ] && procd_append_param command --tun-peer "$tun_peer"

	if [ "$ipv4_only" -eq 1 ]; then
		procd_append_param command --ipv4-only
	else
		[ -n "$tun_local6" ] && procd_append_param command --tun-local6 "$tun_local6"
		[ -n "$tun_peer6" ] && procd_append_param command --tun-peer6 "$tun_peer6"
	fi

	[ -n "$handshake_file" ] && procd_append_param command --handshake-packet "$handshake_file"

	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	
	# Open firewall hole
	procd_open_trigger
	procd_add_reload_trigger "firewall"
	procd_close_trigger

	procd_close_instance
}

start_service() {
	config_load phantun
	config_foreach start_instance phantun
}

service_triggers() {
	procd_add_reload_trigger "phantun"
}
