# luci-app-udp2raw/root/etc/uci-defaults/40_luci-udp2raw
#!/bin/sh
# Copyright (C) 2024 iHub-2020
#
# UCI defaults script for luci-app-udp2raw
# Sets up initial configuration with safe defaults for OpenWrt
#
# Version: 1.0.1
# Last Updated: 2026-01-09

# 只在首次安装或配置不存在时运行
uci -q get udp2raw.general >/dev/null || {
	# 创建全局配置（使用 OpenWrt 安全默认值）
	uci batch <<-EOF
		set udp2raw.general=general
		set udp2raw.general.enabled='0'
		set udp2raw.general.keep_rule='1'
		set udp2raw.general.wait_lock='1'
		set udp2raw.general.retry_on_error='1'
		set udp2raw.general.log_level='4'
		commit udp2raw
	EOF
}

# 确保 rpcd ACL 权限正确
[ -f /usr/share/rpcd/acl.d/luci-app-udp2raw.json ] || {
	mkdir -p /usr/share/rpcd/acl.d
}

# 重启 rpcd 使 ACL 生效
[ -x /etc/init.d/rpcd ] && /etc/init.d/rpcd reload 2>/dev/null

exit 0