# luci-app-udp2raw/root/etc/uci-defaults/40_luci-udp2raw
{
	"luci-app-udp2raw": {
		"description": "Grant access to udp2raw-tunnel configuration",
		"read": {
			"uci": [ "udp2raw" ],
			"ubus": {
				"service": [ "list" ],
				"uci": [ "get", "state" ]
			}
		},
		"write": {
			"uci": [ "udp2raw" ],
			"ubus": {
				"service": [ "list", "start", "stop", "restart", "reload" ],
				"uci": [ "add", "set", "delete", "commit", "apply", "rollback" ]
			}
		}
	}
}