# luci-app-udp2raw/root/etc/uci-defaults/40_luci-udp2raw
#!/bin/sh

uci -q get udp2raw.general || {
	uci set udp2raw.general=general
	uci set udp2raw.general.enabled=0
	uci set udp2raw.general.daemon_user=root
	uci commit udp2raw
}

exit 0
