# luci-app-udp2raw/root/etc/init.d/udp2raw
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2026 Reyanmatic
# Init script for udp2raw
#

START=90
STOP=10
USE_PROCD=1

PROG=/usr/bin/udp2raw

start_service() {
	local enabled daemon_user
	
	config_load udp2raw
	config_get enabled general enabled 0
	[ "$enabled" = "1" ] || return 0
	
	config_get daemon_user general daemon_user "root"
	config_foreach start_tunnel tunnel
}

start_tunnel() {
	local section="$1"
	local mode local_addr local_port remote_addr remote_port
	local raw_mode key cipher_mode auth_mode
	local auto_rule keep_rule seq_mode lower_level
	local source_ip source_port log_level extra_options
	
	config_get mode "$section" mode
	config_get local_addr "$section" local_addr
	config_get local_port "$section" local_port
	config_get remote_addr "$section" remote_addr
	config_get remote_port "$section" remote_port
	config_get raw_mode "$section" raw_mode "faketcp"
	config_get key "$section" key "secret"
	config_get cipher_mode "$section" cipher_mode "aes128cbc"
	config_get auth_mode "$section" auth_mode "md5"
	config_get_bool auto_rule "$section" auto_rule 1
	config_get_bool keep_rule "$section" keep_rule 0
	config_get seq_mode "$section" seq_mode "3"
	config_get lower_level "$section" lower_level
	config_get source_ip "$section" source_ip
	config_get source_port "$section" source_port
	config_get log_level "$section" log_level "4"
	config_get extra_options "$section" extra_options
	
	# 验证必填项
	[ -z "$mode" -o -z "$local_addr" -o -z "$local_port" ] && return 1
	[ -z "$remote_addr" -o -z "$remote_port" ] && return 1
	
	procd_open_instance "$section"
	procd_set_param command $PROG
	
	# 模式
	case "$mode" in
		client) procd_append_param command -c ;;
		server) procd_append_param command -s ;;
		*) return 1 ;;
	esac
	
	# 地址
	procd_append_param command -l "${local_addr}:${local_port}"
	procd_append_param command -r "${remote_addr}:${remote_port}"
	
	# 核心参数
	procd_append_param command -k "$key"
	procd_append_param command --raw-mode "$raw_mode"
	procd_append_param command --cipher-mode "$cipher_mode"
	procd_append_param command --auth-mode "$auth_mode"
	
	# iptables 规则
	[ "$auto_rule" = "1" ] && procd_append_param command -a
	[ "$keep_rule" = "1" ] && procd_append_param command --keep-rule
	
	# 高级选项
	[ -n "$seq_mode" ] && procd_append_param command --seq-mode "$seq_mode"
	[ -n "$lower_level" ] && procd_append_param command --lower-level "$lower_level"
	[ -n "$source_ip" ] && procd_append_param command --source-ip "$source_ip"
	[ -n "$source_port" ] && procd_append_param command --source-port "$source_port"
	[ -n "$log_level" ] && procd_append_param command --log-level "$log_level"
	
	# 额外参数
	for opt in $extra_options; do
		procd_append_param command $opt
	done
	
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param user "${daemon_user}"
	
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "udp2raw"
}
