# luci-app-phantun/root/etc/uci-defaults/40_luci-app-phantun
#!/bin/sh

# 1. Initialize Config (if missing)
if [ ! -f /etc/config/phantun ]; then
	touch /etc/config/phantun
	chmod 644 /etc/config/phantun
	uci set phantun.general=general
	uci set phantun.general.enabled='0'
	uci commit phantun
fi

# 2. Register ucitrack
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@phantun[-1]
	add ucitrack phantun
	set ucitrack.@phantun[-1].init=phantun
	commit ucitrack
EOF

# 3. Reload rpcd for ACLs
if [ -x /etc/init.d/rpcd ]; then
    /etc/init.d/rpcd reload 2>/dev/null || true
fi

# 4. Clear LuCI Cache (Safe Mode)
rm -rf /tmp/luci-indexcache 2>/dev/null || true
rm -rf /tmp/luci-modulecache 2>/dev/null || true

exit 0
